# lib/cap_recipes/tasks/statsd/statsd.conf
description "statsd"
author      "rrussell"

start on (local-filesystems
    and net-device-up IFACE!=lo)
stop on shutdown

respawn

script
    # We found $HOME is needed. Without it, we ran into problems
    export HOME="/root"

    exec sudo -u nobody /usr/share/statsd/scripts/start
end script

post-start script
    # Make a .pid file for god monitoring of statsd
    PID=`status statsd |egrep -oi '([0-9]+)$' |head -n1`
    echo $PID > /var/run/statsd.pid
    ....
end script
 
post-stop script
   # Remove the .pid file used for god monitoring of statsd
   rm -f /var/run/statsd.pid
end script